# vue/docker/vue/s6/etc/services.d/screenshotservice/run
#!/usr/bin/with-contenv sh

cd /opt/vue/screenshot-service
node screenshot-server.js
