# django/docker/django/s6/etc/services.d/django/run
#!/usr/bin/with-contenv sh

mkdir -p /var/services/run
mkdir -p /var/services/django/log

cd /opt/app
uvicorn main.asgi:application --port 8081 --host 0.0.0.0
